# bilibot@.service
[Unit]
Description=Bilibili Comment Auto-Reply Bot
After=network.target

[Service]
Type=simple
User=%I
WorkingDirectory=/home/%I/bilibili-bot
Environment=PYTHONUNBUFFERED=1
Environment=BILIBILI_BOT_ROOT=/home/%I/bilibili-bot
ExecStart=/usr/bin/python3 -m bilibili_bot
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bilibot

[Install]
WantedBy=multi-user.target
